# bin/veritech/veritech.service
[Unit]
Description=Veritech
After=network.target

[Service]
ExecStartPre=-/usr/bin/pkill -e -f '^/firecracker --id'
ExecStart=/usr/local/bin/veritech
ExecStopPost=-/usr/bin/pkill -e -f '^/firecracker --id'
LimitNOFILE=1048576

Type=exec
Restart=always

[Install]
WantedBy=default.target
RequiredBy=network.target