# bin/veritech/veritech.service
[Unit]
Description=Veritech Server
After=network.target

[Service]
ExecStart=/usr/local/bin/veritech --cyclone-local-firecracker --cyclone-pool-size 100 --nats-url $SI_NATS_URL
Type=exec
Restart=always


[Install]
WantedBy=default.target
RequiredBy=network.target
