# components/si-veritech/supervisor/dockerd.conf
[program:dockerd]
command=/usr/bin/dockerd
autostart=true
autorestart=true
